# PerimeterControl-supervisor.service
[Unit]
Description=PerimeterControl Supervisor
Documentation=https://github.com/avgerion/ha-perimeter-control
After=network.target PerimeterControl.service
Wants=PerimeterControl.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/PerimeterControl
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/opt/PerimeterControl"

# Supervisor is run as root so it can manage systemd units via subprocess.
# Future: drop to unprivileged user + sudoers + capability bounding set.
ExecStart=/opt/PerimeterControl/venv/bin/python3 -m supervisor \
    --config /mnt/PerimeterControl/conf \
    --state  /opt/PerimeterControl/state  \
    --port   8080                 \
    --log-level INFO

Restart=on-failure
RestartSec=5

# Give the process time to flush state on shutdown
TimeoutStopSec=30

# Harden the service
PrivateTmp=yes
ProtectSystem=full
ReadWritePaths=/opt/PerimeterControl/state /mnt/PerimeterControl/conf /var/log/PerimeterControl

StandardOutput=journal
StandardError=journal
SyslogIdentifier=PerimeterControl-supervisor

[Install]
WantedBy=multi-user.target
