# system_services/isolator.service
[Unit]
Description=Network Isolator - WiFi AP with Traffic Isolation
Documentation=https://github.com/your-repo/NetworkIsolator
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Configure static IP on whichever interface is isolated by config
ExecStartPre=/opt/PerimeterControl/venv/bin/python3 /opt/PerimeterControl/scripts/network-topology.py --config /mnt/PerimeterControl/conf/isolator.conf.yaml prepare-interface

# Generate and apply configurations
ExecStart=/opt/PerimeterControl/venv/bin/python3 /opt/PerimeterControl/scripts/apply-rules.py --config /mnt/PerimeterControl/conf/isolator.conf.yaml --output-dir /etc/PerimeterControl

# Reload configurations on changes
ExecReload=/opt/PerimeterControl/venv/bin/python3 /opt/PerimeterControl/scripts/apply-rules.py --config /mnt/PerimeterControl/conf/isolator.conf.yaml --output-dir /etc/PerimeterControl

# Cleanup on stop
ExecStop=/usr/sbin/nft flush ruleset
ExecStop=/opt/PerimeterControl/venv/bin/python3 /opt/PerimeterControl/scripts/network-topology.py --config /mnt/PerimeterControl/conf/isolator.conf.yaml cleanup-interface

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=isolator

# Security
ProtectSystem=strict
ReadWritePaths=/etc/PerimeterControl /etc/hostapd /etc/dnsmasq.d
ProtectHome=yes
NoNewPrivileges=false
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
