# system_services/isolator-dashboard.service
[Unit]
Description=Network Isolator Web Dashboard
After=network.target isolator.service
Wants=isolator.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/PerimeterControl/web
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/PerimeterControl/venv/bin/python3 /opt/PerimeterControl/web/dashboard.py
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/PerimeterControl /mnt/PerimeterControl

[Install]
WantedBy=multi-user.target
